# src/config/crontab
0 */3 * * * /bin/python inspect_nodes.py